# templates/etc/cron.d/database_daily_backup.j2
0 1 * * * root /usr/bin/sh {{system.directories.general_scripts}}database_daily_backup.sh > /dev/null