# templates/etc/cron.d/sample_app.j2
0 5 * * * root /usr/bin/php {{ app.directories.web.main }}console.php email diario > /dev/null
0 6 * * 0 root /usr/bin/php {{ app.directories.web.main }}console.php email semanal > /dev/null